# Library/LaunchAgents/com.ivanzoid.ssh-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ivanzoid.ssh-tunnel</string>
  <key>OnDemand</key>
  <false/>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/ssh</string>
    <string>-N</string>
    <string>-o</string>
    <string>StrictHostKeyChecking=accept-new</string>
    <string>-o</string>
    <string>ControlMaster=no</string>
    <string>-o</string>
    <string>ControlPath=none</string>
    <string>-o</string>
    <string>ExitOnForwardFailure=yes</string>
    <string>-o</string>
    <string>ServerAliveInterval=15</string>
    <string>-o</string>
    <string>ServerAliveCountMax=3</string>
    <string>-o</string>
    <string>TCPKeepAlive=no</string>
    <string>-D</string>
    <string>8000</string>
    <string>ivan@vpn.zoid.cc</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/ivan/Library/Logs/ssh-tunnel.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ivan/Library/Logs/ssh-tunnel.err.log</string>
</dict>
</plist>
